# pure-ftpd.socket
[Unit]
Description=Pure-FTPd Socket

[Socket]
ListenStream=21
Accept=yes

[Install]
WantedBy=sockets.target
